# ops/backup/systemd/ligate-backup@.service
[Unit]
Description=Ligate Chain RocksDB backup to GCS (tier %i)
# Don't pile up backups if the previous one is still running (large
# state + slow network = backup runs longer than the timer interval).
After=ligate-node.service

[Service]
Type=oneshot
User=ligate
EnvironmentFile=/etc/ligate/backup.env
# Template unit: the tier (hourly | daily | weekly) is passed as the
# instance name `%i`, so a single .service file backs all three
# timers without per-tier drop-ins.
#
#   ligate-backup@hourly.service   →   TIER=hourly  (gs:// .../hourly/<HH>/)
#   ligate-backup@daily.service    →   TIER=daily   (gs:// .../daily/<YYYYMMDD>/)
#   ligate-backup@weekly.service   →   TIER=weekly  (gs:// .../weekly/<YYYY-Www>/)
#
# `backup.env` provides at minimum:
#   GCS_BUCKET=ligate-devnet-1-backups
#   LIGATE_STORAGE_DIR=...   (override; default is /var/lib/ligate/rollup)
# It does NOT set TIER any more (that comes from %i below) and would
# be silently overwritten if also in the envfile.
Environment=TIER=%i
ExecStart=/opt/ligate/scripts/backup-rocksdb.sh
# Two metrics matter on backup ops: success / failure + duration.
# Both surface in `journalctl -u ligate-backup@*` and the next
# observability pass (#268) will wire a Prometheus textfile exporter
# to scrape them.
StandardOutput=journal
StandardError=journal
