# ops/backup/systemd/ligate-backup-daily.timer
[Unit]
Description=Ligate Chain RocksDB daily backup (keep-7 tier)
Requires=ligate-backup.service

[Timer]
# 03:30 UTC nightly. Tier flag set via drop-in below.
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=ligate-backup.service

[Install]
WantedBy=timers.target
