# ops/backup/systemd/ligate-backup-hourly.timer
[Unit]
Description=Ligate Chain RocksDB hourly backup
Requires=ligate-backup@hourly.service

[Timer]
# 15-minute cadence per the issue scope; hourly tier files them under
# `gs://$GCS_BUCKET/<host>/hourly/<timestamp>/`. Local rotation keeps
# last 24 hourly snapshots on disk.
OnCalendar=*:0/15
# Catch up missed runs if the node was down at the scheduled tick.
Persistent=true
# Anchor to the boundary minute (so timer fires at :00 :15 :30 :45,
# not at process-start time + 15m).
AccuracySec=1s
Unit=ligate-backup@hourly.service

[Install]
WantedBy=timers.target
